# deployment/gunicorn_portfolio.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn_portfolio.sock

[Install]
WantedBy=sockets.target
